# docker/debian-dev/rootfs/etc/profile.d/zzz-custom.sh
#!/bin/sh -xe

export LESS='-S -M -R -i -#10'

alias ll='ls -al --color'
